# debian/adpi-utils/lib/udev/rules.d/99-adpi-iio.rules
ACTION=="add", SUBSYSTEM=="iio", ATTRS{name}=="ad7794", RUN+="/usr/sbin/adpictl --iio %m set gain 1"
